# supervisord-nginx.conf
[program:nginx]
command=/opt/nginx/sbin/nginx
