# supervisord-rails.conf
[program:rails]
command=rails server
directory=/app
