# supervisord-nginx.conf
[program:nginx]
command=/opt/nginx/sbin/nginx
numprocs=1
autostart=true
autorestart=true
